# systemd/uploader.timer
[Unit]
Description=Mothpi Pictures Uploader Timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
